# systemd/import-cityleague-result-job_enqueue.service
[Unit]
Description=import-cityleague-result-job_enqueue
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/mkr wrap --name import-cityleague-result-job_enqueue --auto-close -- /home/ubuntu/vsrecorder/import-cityleague-result-job/bin/enqueue
WorkingDirectory=/home/ubuntu/vsrecorder/import-cityleague-result-job

# 失敗時リトライ
Restart=on-failure
RestartSec=900
StartLimitInterval=5400
StartLimitBurst=5
