# systemd/import-cityleague-result-job_enqueue.service
[Unit]
Description=import-cityleague-result-job_enqueue
After=network.target

[Service]
Type=oneshot
ExecStart=/bin/bash -lc '/usr/bin/mkr wrap --name import-cityleague-result-job_enqueue --auto-close -- /home/ubuntu/vsrecorder/import-cityleague-result-job/bin/enqueue && ts=$(date +%%s); printf "import-cityleague-result-job.enqueue.last_run_time\t%%s\t%%s\n" "$ts" "$ts" | mkr throw --service monolith'
WorkingDirectory=/home/ubuntu/vsrecorder/import-cityleague-result-job

# 失敗時リトライ
Restart=on-failure
RestartSec=900
StartLimitInterval=5400
StartLimitBurst=5
