# systemd/import-cityleague-result-job_enqueue.timer
[Unit]
Description=import-cityleague-result-job_enqueue

[Timer]
OnCalendar=*-*-* 14:15:00
Persistent=true

[Install]
WantedBy=timers.target